# deploy/systemd/octopos-telegram.service.template
[Unit]
Description=octopOS Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
__SERVICE_USER_LINE__
WorkingDirectory=__WORKDIR__
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=__WORKDIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__PYTHON_BIN__ -m src.interfaces.cli.main telegram
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=__WANTED_BY__
